# external/vl53l1_daemon/vl53l1-daemon.service
[Unit]
Description=VL53L1 ToF sensor daemon
SourcePath=/usr/bin/vl53l1_daemon
Before=vic-robot.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/vl53l1_daemon
