# external/fake-hwclock/fake-hwclock-cmd.socket
[Unit]
Description=Fake Hardware Clock socket

[Socket]
ListenFIFO=/run/fake-hwclock-cmd

[Install]
WantedBy=sockets.target
