# external/fake-hwclock/fake-hwclock.service
# Fake hwclock service
# No default dependencies to enable very early execution.  Depends on mount-data
# Manually add shutdown.target dependency because no default.
# pulls in time-sync.target (see man systemd.special)
# Runs before the other time sync daemons, if they are being started

[Unit]
Description=Fake Hardware Clock
DefaultDependencies=No
Requires=mount-data.service
After=mount-data.service
Wants=time-sync.target
Conflicts=shutdown.target
Before=shutdown.target time-sync.target ntpdate.service ntpd.service chronyd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/fake-hwclock load
ExecStop=/sbin/fake-hwclock save

[Install]
WantedBy=multi-user.target
Also=fake-hwclock-tick.timer
