# external/rebooter/rebooter.service
# Service to periodically reboot the robot

[Unit]
Description=Robot Rebooter
After=fake-hwclock.service
Requires=fake-hwclock.service

[Service]
Type=simple
ExecStart=/sbin/rebooter.py

