# oe/meta-anki/recipes/fault-code/files/fault-code.socket
[Unit]
Description=Fault Code socket

[Socket]
ListenFIFO=/run/fault_code

[Install]
WantedBy=sockets.target
