# poky/meta-anki/recipes/anki-crash-log/files/anki-crash-log.socket
[Unit]
Description=Anki Crash Log Socket

[Socket]
ListenFIFO=/run/anki-crash-log

[Install]
WantedBy=sockets.target
