# poky/meta-anki/recipes/spacedog/files/spacedog.service
[Unit]
Description=@DESCRIPTION@

[Service]
ExecStart=@BINDIR@/spacedog.sh

[Install]
WantedBy=multi-user.target
